# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
pidfile=/var/run/supervisord.pid
user=root
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; =====================
; 通用模板（示例，默认注释）
; =====================
[program:python-worker]
command=python3 app.py
;directory=/app
user=root
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/pworker.out.log
stderr_logfile=/var/log/supervisor/pworker.err.log

; =====================
; 通用模板（示例，默认注释）
; =====================
[program:nodejs-worker]
command=node app.js
;directory=/app
user=root
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/nworker.out.log
stderr_logfile=/var/log/supervisor/nworker.err.log

; =====================
; 允许 include 扩展
; =====================
[include]
files = /etc/supervisor/conf.d/*.conf
